# target/linux/x86/base-files/etc/init.d/mdadm_hack
#!/bin/sh /etc/rc.common
# work around mdadm create /dev/md* fail
# Copyright (C) 2022 jjm2473

START=12

boot() {
	[ -f /proc/mdstat ] || modprobe md_mod
	mount --bind /dev/full /sys/module/md_mod/parameters/new_array
}
