# target/linux/x86/base-files/etc/uci-defaults/kmods-unload-x86
#!/bin/sh

[ -f /etc/config/kmods ] || exit 0

uci -q batch <<-EOF >/dev/null
	delete kmods.amdgpu
	set kmods.amdgpu=auto_unload
	set kmods.amdgpu.function=amdgpu
	delete kmods.i915
	set kmods.i915=auto_unload
	set kmods.i915.function=i915
	delete kmods.radeon
	set kmods.radeon=auto_unload
	set kmods.radeon.function=radeon

	delete kmods.mwifiex
	set kmods.mwifiex=auto_unload
	set kmods.mwifiex.function=mwifiex_pcie
	set kmods.mwifiex.support=mwifiex
	delete kmods.ath10k
	set kmods.ath10k=auto_unload
	set kmods.ath10k.function=ath10k_pci
	set kmods.ath10k.support=ath10k_core

	delete kmods.qed
	set kmods.qed=auto_unload
	set kmods.qed.function="qede qedf"
	set kmods.qed.support=qed

	delete kmods.x86_others
	set kmods.x86_others=auto_unload
	set kmods.x86_others.function="atlantic bnx2x mlx5_core mlx4_core i40e mpt3sas ixgbe r8126"

	commit kmods
EOF
