# target/linux/realtek/base-files/etc/uci-defaults/11_fstab_automount
#!/bin/sh

if [ ! -f /.recovery_mode ]; then
    uci set fstab.@global[0].anon_mount=1
    uci commit fstab
fi

exit 0
