# target/linux/realtek/base-files/etc/uci-defaults/31_upgrade21
#!/bin/sh

if ! [ -f /etc/uci-defaults/99_theme ]; then
    if [ -d /www/luci-static/argon ]; then
        uci set luci.main.mediaurlbase=/luci-static/argon
        uci commit luci
    fi
fi
