# package/system/fstools/files/fstab.init
#!/bin/sh /etc/rc.common
# (C) 2013 openwrt.org

START=40

boot() {
	/sbin/block mount
}

start() {
	return 0
}

restart() {
	/sbin/block mount
}

stop() {
	/sbin/block umount
}
