# package/istoreos-files/files/etc/hotplug.d/block/09-blockmount
if [ "$ACTION" = "add" ]; then
    [ "`uci -q get fstab.@global[0].anon_mount`" = "1" ] && /usr/libexec/blockmount.sh
fi