# package/base-files/files/etc/hotplug.d/block/11-umount
if [ "$ACTION" = "remove" ]; then
    umount "/dev/$DEVNAME"
fi
