# package/boot/rtk-init/files/printk.init
#!/bin/sh /etc/rc.common
# 2020 jjm2473 porting openwrt 19.07+ to zidoo z9s

START=99

boot() {
	if [ -n "$ANDROID" ]; then
		# tune printk level after success boot, avoid noise in console
		echo '4 5 1 7' > /proc/sys/kernel/printk
	fi
}
